# cli/hightower-stun.service
[Unit]
Description=Hightower STUN Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ht run stun
Restart=on-failure
RestartSec=5
User=hightower
Group=hightower

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/hightower

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hightower-stun

[Install]
WantedBy=multi-user.target
